# deploy/poly-dashboard.service
[Unit]
Description=Poly Autobetting Dashboard API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/poly-autobetting
EnvironmentFile=/home/ubuntu/poly-autobetting/.env
ExecStart=/home/ubuntu/poly-autobetting/venv/bin/uvicorn dashboard.server:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=poly-dashboard

[Install]
WantedBy=multi-user.target

